# deploy/systemd/2salti.service.d/override.conf
[Service]
ExecReload=/bin/kill -s HUP $MAINPID
