# deploy/systemd/dev/2salti-dev-ocrworker.service
[Unit]
Description=2salti OCR worker (coda referti) - DEV
After=network.target

[Service]
User=alberto
Group=www-data
WorkingDirectory=/opt/2salti-dev
ExecStart=/opt/2salti-dev/.venv/bin/python manage.py ocr_worker
Restart=always
RestartSec=5
# SIGTERM: il worker finisce il job in corso (~80s con Gemini) e poi esce.
TimeoutStopSec=150

[Install]
WantedBy=multi-user.target
